# deploy/testnet/solen-node.service
[Unit]
Description=Solen Testnet Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=solen
Group=solen
WorkingDirectory=/opt/solen
ExecStart=/opt/solen/bin/solen-node \
    --network testnet \
    --genesis /opt/solen/config/genesis.json \
    --data-dir /opt/solen/data/testnet \
    --bootstrap /ip4/seed1.solenchain.com/tcp/40333 \
    --bootstrap /ip4/seed2.solenchain.com/tcp/40333
Restart=always
RestartSec=5
LimitNOFILE=65536

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
